# var/lib/dpkg/info/ctxusb.list
/.
/etc
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/85-ica-usb.rules
/etc/init.d
/etc/init.d/ctxusbd
/usr
/usr/share
/usr/share/doc
/usr/share/doc/ctxusb
/usr/share/doc/ctxusb/copyright
/usr/share/doc/ctxusb/changelog.gz
/opt
/opt/Citrix
/opt/Citrix/ICAClient
/opt/Citrix/ICAClient/ctxusbd
/opt/Citrix/ICAClient/ctxusb
/opt/Citrix/ICAClient/ctx_usb_isactive
/opt/Citrix/ICAClient/usb.conf
/opt/Citrix/ICAClient/VDGUSB.DLL
