# var/lib/dpkg/info/openvpn.list
/.
/usr
/usr/share
/usr/share/openvpn
/usr/share/openvpn/verify-cn
/usr/share/man/man8
/usr/share/man/man8/openvpn.8.gz
/usr/share/doc
/usr/share/doc/openvpn
/usr/share/doc/openvpn/COPYRIGHT.GPL.gz
/usr/share/doc/openvpn/README
/usr/share/doc/openvpn/README.Debian.gz
/usr/share/doc/openvpn/AUTHORS
/usr/share/doc/openvpn/PORTS
/usr/share/doc/openvpn/README.IPv6
/usr/share/doc/openvpn/NEWS.Debian.gz
/usr/share/doc/openvpn/management-notes.txt.gz
/usr/share/doc/openvpn/README.down-root
/usr/share/doc/openvpn/README.polarssl
/usr/share/doc/openvpn/changelog.Debian.gz
/usr/share/doc/openvpn/changelog.gz
/usr/share/doc/openvpn/COPYING.gz
/usr/share/doc/openvpn/copyright
/usr/share/doc/openvpn/README.auth-pam
/usr/share/doc/openvpn/examples
/usr/share/doc/openvpn/examples/sample-keys
/usr/share/doc/openvpn/examples/sample-keys/README
/usr/share/doc/openvpn/examples/sample-keys/dh1024.pem
/usr/share/doc/openvpn/examples/sample-keys/pass.crt
/usr/share/doc/openvpn/examples/sample-keys/pkcs12.p12
/usr/share/doc/openvpn/examples/sample-keys/client.key
/usr/share/doc/openvpn/examples/sample-keys/ca.crt
/usr/share/doc/openvpn/examples/sample-keys/pass.key
/usr/share/doc/openvpn/examples/sample-keys/client.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/server.key
/usr/share/doc/openvpn/examples/sample-keys/server.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/ca.key
/usr/share/doc/openvpn/examples/sample-config-files
/usr/share/doc/openvpn/examples/sample-config-files/home.up
/usr/share/doc/openvpn/examples/sample-config-files/firewall.sh
/usr/share/doc/openvpn/examples/sample-config-files/static-office.conf
/usr/share/doc/openvpn/examples/sample-config-files/README
/usr/share/doc/openvpn/examples/sample-config-files/tls-office.conf
/usr/share/doc/openvpn/examples/sample-config-files/server.conf.gz
/usr/share/doc/openvpn/examples/sample-config-files/static-home.conf
/usr/share/doc/openvpn/examples/sample-config-files/openvpn-startup.sh
/usr/share/doc/openvpn/examples/sample-config-files/xinetd-server-config
/usr/share/doc/openvpn/examples/sample-config-files/client.conf
/usr/share/doc/openvpn/examples/sample-config-files/office.up
/usr/share/doc/openvpn/examples/sample-config-files/tls-home.conf
/usr/share/doc/openvpn/examples/sample-config-files/xinetd-client-config
/usr/share/doc/openvpn/examples/sample-config-files/openvpn-shutdown.sh
/usr/share/doc/openvpn/examples/sample-config-files/loopback-client
/usr/share/doc/openvpn/examples/sample-config-files/loopback-server
/usr/share/doc/openvpn/examples/sample-scripts
/usr/share/doc/openvpn/examples/sample-scripts/ucn.pl
/usr/share/doc/openvpn/examples/sample-scripts/verify-cn
/usr/share/doc/openvpn/examples/sample-scripts/auth-pam.pl
/usr/share/doc/openvpn/examples/sample-scripts/bridge-stop
/usr/share/doc/openvpn/examples/sample-scripts/bridge-start
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/openvpn.conf
/usr/lib/openvpn
/usr/lib/openvpn/openvpn-plugin-auth-pam.so
/usr/lib/openvpn/openvpn-plugin-down-root.so
/usr/sbin
/usr/sbin/openvpn
/usr/include
/usr/include/openvpn
/usr/include/openvpn/openvpn-plugin.h
/etc
/etc/bash_completion.d
/etc/bash_completion.d/openvpn
/etc/network
/etc/network/if-down.d
/etc/network/if-down.d/openvpn
/etc/network/if-up.d
/etc/network/if-up.d/openvpn
/etc/default
/etc/default/openvpn
/etc/openvpn
/etc/openvpn/update-resolv-conf
/etc/init.d
/etc/init.d/openvpn
/lib
/lib/systemd
/lib/systemd/system-generators
/lib/systemd/system-generators/openvpn-generator
/lib/systemd/system
/lib/systemd/system/openvpn@.service
/lib/systemd/system/openvpn.service
