# var/lib/dpkg/info/udisks.list
/.
/usr
/usr/lib
/usr/lib/udisks
/usr/lib/udisks/udisks-helper-create-partition-table
/usr/lib/udisks/udisks-helper-mkfs
/usr/lib/udisks/udisks-helper-fstab-mounter
/usr/lib/udisks/udisks-helper-create-partition
/usr/lib/udisks/udisks-helper-change-luks-password
/usr/lib/udisks/udisks-helper-drive-poll
/usr/lib/udisks/udisks-helper-linux-md-remove-component
/usr/lib/udisks/udisks-helper-drive-detach
/usr/lib/udisks/udisks-helper-change-filesystem-label
/usr/lib/udisks/udisks-helper-drive-benchmark
/usr/lib/udisks/udisks-helper-modify-partition
/usr/lib/udisks/udisks-helper-ata-smart-selftest
/usr/lib/udisks/udisks-helper-ata-smart-collect
/usr/lib/udisks/udisks-helper-delete-partition
/usr/lib/udisks/udisks-helper-linux-md-check
/usr/lib/udisks/udisks-helper-mdadm-expand
/usr/lib/udisks/udisks-daemon
/usr/bin
/usr/bin/udisks
/usr/bin/udisks-tcp-bridge
/usr/share
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.udisks.policy
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.UDisks.service
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.UDisks.Adapter.xml
/usr/share/dbus-1/interfaces/org.freedesktop.UDisks.Expander.xml
/usr/share/dbus-1/interfaces/org.freedesktop.UDisks.Device.xml
/usr/share/dbus-1/interfaces/org.freedesktop.UDisks.xml
/usr/share/dbus-1/interfaces/org.freedesktop.UDisks.Port.xml
/usr/share/locale
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/udisks.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/udisks.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/udisks.mo
/usr/share/doc
/usr/share/doc/udisks
/usr/share/doc/udisks/changelog.Debian.armhf.gz
/usr/share/doc/udisks/README
/usr/share/doc/udisks/NEWS.gz
/usr/share/doc/udisks/AUTHORS
/usr/share/doc/udisks/changelog.Debian.gz
/usr/share/doc/udisks/copyright
/usr/share/man/man8
/usr/share/man/man8/udisks-daemon.8.gz
/usr/share/man/man7
/usr/share/man/man7/udisks.7.gz
/usr/share/man/man1
/usr/share/man/man1/udisks-tcp-bridge.1.gz
/usr/share/man/man1/udisks.1.gz
/usr/share/pkgconfig
/usr/share/pkgconfig/udisks.pc
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/udisks.service
/lib/udev
/lib/udev/udisks-lvm-pv-export
/lib/udev/udisks-dm-export
/lib/udev/rules.d
/lib/udev/rules.d/80-udisks.rules
/lib/udev/udisks-probe-sas-expander
/lib/udev/udisks-part-id
/lib/udev/udisks-probe-ata-smart
/var
/var/lib
/var/lib/udisks
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.UDisks.conf
/etc/avahi
/etc/avahi/services
/etc/avahi/services/udisks.service
/sbin
/sbin/umount.udisks
